# profile.d/apt_ignore_files-lzfuie9tjzt6kqjh15zi2gke7.sh
#! /bin/false
# v2025.266

# Sync this with /etc/environment where the same variable is set.
export APT_CONFIG=/etc/apt/apt.conf.d\
/50ignore_upstream-hdq2kzpvyehnsea0gtcr6vflu
