# cron.daily/lifetime-writes
#! /bin/sh
logger updating lifetime-writes
exec sh /var/lib/lifetime-writes/bytes_written-update.sh > /dev/null
