# cron.daily/logrotate-localsite-m2fv383qa8qhrx8udayxe753a
#! /bin/sh
exec /etc/logrotate-localsite/refresh-fixed-upstream-snippets.sh 2> /dev/null
