# cron.d/syslog_ring-8tuapunktw7mm9jgnhtri2zhv
#     m  h dom mon dow user command
7-59/30  *   *   *   * root sh -l /etc/cron.daily/save_logs-4vzntsgfwksyilhv27klxte5l
